# crontab
# START CRON JOB
*/2 * * * * /usr/local/bin/python /app/dollar.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB